# base/conf/supervisor.conf
[unix_http_server]
file=/opt/run/supervisor.sock   ; (the path to the socket file)
chmod=0777
[supervisord]
logfile=/opt/run/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB        ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10           ; (num of main logfile rotation backups;default 10)
loglevel=info                ; (log level;default info; others: debug,warn,trace)
stderr_logfile=/opt/run/error.log;
pidfile=/opt/run/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                  ; (min. avail startup file descriptors;default 1024)
minprocs=200                 ; (min. avail process descriptors;default 200)
childlogdir=/opt/run/ ;
#user=root;

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///opt/run/supervisor.sock ; use a unix:// URL  for a unix socket


[program:m-game-app]

process_name = m-game-app-%(process_num)s

command = /opt/python2.7/bin/python2.7 /opt/app/python/m-game-platform/app/manage.py runfcgi method=prefork host=127.0.0.1 port=%(process_num)s

autostart = true
numprocs_start = 20000
numprocs= 20


[program:m-game-api]

process_name = m-game-api-%(process_num)s

command = /opt/python2.7/bin/python2.7 /opt/app/python/m-game-platform/api/main.py --port=%(process_num)s --logging=error --log_file_prefix=/opt/logs/tornado/game/%(program_name)s-%(process_num)s.log

autostart = true
numprocs_start = 10003
numprocs= 200



